# {{cookiecutter.project_slug}}/supervisord.conf
[program:{{cookiecutter.project_slug}}]
command=gunicorn -w 4 -b 0.0.0.0:5000 {{cookiecutter.project_slug}}:app
directory= ./
autostart=true
autorestart=true 

[supervisord]
logfile=log/supervisord.log                    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=info                      ; info, debug, warn, trace
pidfile=log/supervisord.pid                    ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors